# bin/rc.local
#Add this line to /etc/rc.local
/home/pi/CODE/pantilt-cam/bin/www --start-camera &
